# utils/jslisten.service.in
[Unit]
Description=jslisten (JoyStick Listening) daemon
After=network.target auditd.service

[Service]
Type=idle
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/#USER#/.Xauthority"
ExecStart=/bin/sh -ec "exec #BINDIR#/jslisten"
ExecReload=/bin/kill -9 $MAINPID
ExecStop=/bin/kill -9 $MAINPID
KillMode=process
Restart=on-failure
User=root

[Install]
WantedBy=multi-user.target
